# vagrant/etc/supervisor/conf.d/bower.conf
[program:bower]
command=/usr/bin/bower update
user=vagrant
directory=/pycroft
startretries=0
autorestart=false
priority=900
environment=HOME="/home/vagrant"
