# vagrant/etc/supervisor/conf.d/postgresql.conf
[program:postgresql]
command=/usr/local/bin/postgresql.sh
autorestart=true
